# agent/systemd/auto-facebook-agent-stack.target
[Unit]
Description=fb.autonow.vn agent — browser stack (Xvfb + x11vnc + websockify)
Documentation=https://fb.autonow.vn
Wants=auto-facebook-agent-xvfb.service auto-facebook-agent-x11vnc.service auto-facebook-agent-websockify.service

[Install]
WantedBy=multi-user.target
